# debian/hildon-control-panel-personalisation.init
#!/sbin/openrc-run
### BEGIN INIT INFO
# Provides: hildon-control-panel-personalisation
# Required-Start:
# Required-Stop:
# Default-Start: 2 3 4 5
# Default-Stop:
### END INIT INFO

description="Hildon Control Panel Personalisation"
name="hildon-control-panel-personalisation"

depend() {
        before xorg xsession
}

start() {
    ebegin "Starting hildon-control-panel-personalisation"
	if [ ! -e /etc/hildon/theme/index.theme ]; then
        ebegin "Setting default theme"
        /usr/bin/personalisation /usr/share/themes/alpha;
	fi
}
